# systemd/.config/systemd/user/logitech-g300s-config.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/logitech-g300s-configuration.sh
